# multimedia/couchpotato/files/couchpotato.in
#!/bin/sh
#
# PROVIDE: couchpotato
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# couchpotato_enable:	Set to NO by default. Set it to YES to enable it.
# couchpotato_user:	The user account the CouchPotaato daemon runs as.
#			It uses '%%USERS%%' user by default.
# couchpotato_datadir:	Directory where couchpotato user data is stored.
#			Default: /var/db/couchpotato

. /etc/rc.subr
name="couchpotato"
rcvar="${name}_enable"
load_rc_config ${name}

: ${couchpotato_enable:="NO"}
: ${couchpotato_user:="%%USERS%%"}
: ${couchpotato_datadir:="/var/db/couchpotato"}

pidfile="/var/run/couchpotato/couchpotato.pid"
command="%%PYTHON_CMD%%"
command_args="%%DATADIR%%/CouchPotatoServer/CouchPotato.py --daemon --data_dir ${couchpotato_datadir} --pid_file ${pidfile}"

start_precmd="couchpotato_prestart"
couchpotato_prestart()
{
	export LC_CTYPE="UTF-8.en_US"
	if [ ! -d ${pidfile%/*} ]; then
		install -d -o ${couchpotato_user} ${pidfile%/*}
	fi
	if [ ! -d ${couchpotato_datadir} ]; then
		install -d -o ${couchpotato_user} ${couchpotato_datadir}
	fi
}

run_rc_command "$1"
