# multimedia/sickbeard/files/sickbeard.in
#!/bin/sh
#
# PROVIDE: sickbeard
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# sickbeard_enable:	Set to NO by default. Set it to YES to enable it.
# sickbeard_user:	The user account the SickBeard daemon runs as.
#			It uses '%%USERS%%' by default.
# sickbeard_datadir:	Directory where SickBeard user data is stored.
#			Default: /var/db/sickbeard

. /etc/rc.subr
name="sickbeard"
rcvar="${name}_enable"
load_rc_config ${name}

: ${sickbeard_enable:="NO"}
: ${sickbeard_user:="%%USERS%%"}
: ${sickbeard_datadir:="/var/db/sickbeard"}

pidfile="/var/run/sickbeard/sickbeard.pid"
command="%%PYTHON_CMD%%"
command_args="%%DATADIR%%/Sick-Beard/SickBeard.py --daemon --datadir ${sickbeard_data_dir} --pidfile ${pidfile} --nolaunch"

start_precmd="sickbeard_prestart"
sickbeard_prestart()
{
	if [ -f ${pidfile} ]; then
		rm -f ${pidfile}
		echo "Removing stale pidfile."
	elif [ ! -d ${pidfile%/*} ]; then
		install -d -o ${sickbeard_user} ${pidfile%/*}
	fi
	if [ ! -d ${sickbeard_datadir} ]; then
		install -d -o ${sickbeard_user} ${sickbeard_datadir}
	fi
}

run_rc_command "$1"
