# multimedia/headphones/files/headphones.in
#!/bin/sh
#
# PROVIDE: headphones
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# headphones_enable:	Set to NO by default. Set it to YES to enable it.
# headphones_user:	The user account Headphones daemon runs as.
#			It uses '%%USERS%%' user by default.
# headphones_datadir:	Directory where Headphones user data is stored.
#			Default: /var/db/headphones

. /etc/rc.subr
name="headphones"
rcvar="${name}_enable"
load_rc_config ${name}

: ${headphones_enable:="NO"}
: ${headphones_user:="%%USERS%%"}
: ${headphones_datadir:="/var/db/headphones"}

pidfile="/var/run/headphones/headphones.pid"
command="%%PYTHON_CMD%%"
command_args="%%DATADIR%%/headphones/Headphones.py --daemon --datadir ${headphones_datadir} --pidfile ${pidfile} --nolaunch"

start_precmd="headphones_prestart"
headphones_prestart()
{
	if [ -f $pidfile ]; then
		rm -f $pidfile
		echo "Removing stale pidfile."
	elif [ ! -d ${pidfile%/*} ]; then
		install -d -o ${headphones_user} ${pidfile%/*}
	fi
	if [ ! -d ${headphones_datadir} ]; then
		install -d -o ${headphones_user} ${headphones_datadir}
	fi
}

run_rc_command "$1"
