# multimedia/sickrage/files/sickrage.in
#!/bin/sh
#
# PROVIDE: sickrage
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# sickrage_enable:	Set to NO by default. Set it to YES to enable it.
# sickrage_user:	The user account the SickRage daemon runs as.
#			It uses '%%USERS%%' by default.
# sickrage_datadir:	Directory where SickRage user data lives.
#			Default: /var/db/sickrage

. /etc/rc.subr
name="sickrage"
rcvar="${name}_enable"
load_rc_config ${name}

: ${sickrage_enable:="NO"}
: ${sickrage_user:="%%USERS%%"}
: ${sickrage_datadir:="/var/db/sickrage"}

pidfile="/var/run/sickrage/sickrage.pid"
command="%%PYTHON_CMD%%"
command_args="%%DATADIR%%/SickBeard.py --daemon --datadir ${sickrage_datadir} --pidfile ${pidfile} --nolaunch"

start_precmd="sickrage_prestart"
sickrage_prestart()
{
	export LC_CTYPE="UTF-8.en_US"
	if [ ! -d ${pidfile%/*} ]; then
		install -d -o ${sickrage_user} ${pidfile%/*}
	fi
	if [ ! -d ${sickrage_datadir} ]; then
		install -d -o ${sickrage_user} ${sickrage_datadir}
	fi
}

run_rc_command "$1"
